# deploy/lunda-admin.service
[Unit]
Description=Lunda Admin Panel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lunda_admin
EnvironmentFile=/opt/lunda_admin/.env
ExecStart=/opt/lunda_admin/venv/bin/uvicorn admin.app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
